# scripts/nextcloud.service
[Unit]
Description=Docker - Nextcloud container
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/docker/nextcloud_docker/scripts/start.sh
ExecStop=/docker/nextcloud_docker/scripts/stop.sh

[Install]
WantedBy=default.target
				
